# ansible/playbooks/roles/rpi_init/templates/rc.local.j2
#!/bin/sh -e

# setup rpi's dashboard
cd /home/rpi/U6143_ssd1306/C
sudo make clean
sudo make
sudo ./display &